# roles/matrix/nginx-proxy-health/templates/systemd/matrix-publish-health.service.j2
[Unit]
Description=Publish Matrix Services status to html report

[Service]
Type=oneshot
ExecStart={{ matrix_local_bin_path }}/matrix-publish-health
