# roles/system/maintenance/templates/systemd/custom-journalctl-vacuum.service.j2
[Unit]
Description=Run journalctl vacuuum procedure

[Service]
Type=oneshot
ExecStart=/bin/journalctl --vacuum-time=7d
